# sample/config/simple_ap.service.example
[Unit]
Description=simple ActibityPub Server

[Service]
Type=simple
TimeoutSec=60
StandardOutput=journal
StandardError=journal
SyslogIdentifier=simple_ap
Restart=always
